# bin/ssl-proxy.service
Description=SSL proxy

Wants=network.target
After=syslog.target network-online.target multi-user.target

[Service]
Type=simple
ExecStart=/home/bgf/bin/ssl-proxy.sh
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
